# bin/webpc.service
[Unit]
Description=Web PC Service
After=network.target
 
[Service]
#User=webpc
Type=simple
ExecStart=/opt/webpc/webpc daemon -r
KillMode=control-group
Restart=on-failure
 
[Install]
WantedBy=multi-user.target
